# service/haproxy/run
#!/bin/sh
exec 2>&1
exec /usr/local/sbin/haproxy-systemd-wrapper -p /run/haproxy.pid -f /marathon-lb/haproxy.cfg
